# merlind.service.in
[Unit]
Description=Merlin
After=syslog.target network.target mariadb.service

[Service]
#Starting the merlind with --debug option won't daemonize it.
#By doing this we leverage the rhel7 systemd unit file structure.
ExecStart=/usr/bin/merlind --config @pkgconfdir@/merlin.conf --debug
Restart=always
User=monitor
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
