# apps/reportdata-push.cron
# push logs 02:00 every night
00 02 * * * monitor mon log push
